# ansible/all-in-one/roles/xsnippet/templates/cert-watcher.path.j2
[Path]
PathModified={{ watch_path }}

[Install]
WantedBy=multi-user.target
