# ansible/all-in-one/roles/sockets/files/http.service
[Unit]
Requires=http.socket
After=http.socket

[Service]
ExecStart=/usr/lib/systemd/systemd-socket-proxyd localhost:8080
