# ansible/all-in-one/roles/sockets/files/https.service
[Unit]
Requires=https.socket
After=https.socket

[Service]
ExecStart=/usr/lib/systemd/systemd-socket-proxyd localhost:8443
